# infra/release/cursor-gateway-release-sync.service
[Unit]
Description=Mirror interview reports to Cursor Gateway release
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/joel/cursor-gateway-release/integrations/release/sync_reports.py --limit 7
UMask=0077
