# integrations/hermes/systemd/hermes-gateway.service.d/zz-resource-limits.conf
[Service]
MemoryAccounting=true
MemoryHigh=640M
MemoryMax=768M
MemorySwapMax=256M
TasksMax=128
OOMPolicy=stop
