# integrations/hermes/systemd/hermes-gateway.service.d/zz-resource-limits.conf
[Service]
MemoryAccounting=true
MemoryHigh=500M
MemoryMax=650M
TasksMax=128
OOMPolicy=stop
