# integrations/hermes/systemd/hermes-gateway-telegram2.service.d/zz-resource-limits.conf
[Service]
MemoryAccounting=true
MemoryHigh=384M
MemoryMax=512M
MemorySwapMax=128M
TasksMax=128
OOMPolicy=stop
