# integrations/hermes/systemd/hermes-gateway-telegram2.service.d/zz-resource-limits.conf
[Service]
MemoryAccounting=true
MemoryHigh=430M
MemoryMax=550M
TasksMax=128
OOMPolicy=stop
