# integrations/hermes/systemd/hermes-gateway-telegram2.service.d/zz-strict-route.hermes-ha.conf
# Corrected for hermes-ha profile path. Install as zz-strict-route.conf only after
# verifying preflight against this HERMES_HOME.
# Requires/After keep the watcher wanted while the gateway is running.
# Intentionally no BindsTo — avoid health-blip guard races stopping the gateway.
[Unit]
Requires=hermes-strict-route-guard@telegram2.service
After=hermes-strict-route-guard@telegram2.service

[Service]
Environment=HERMES_HOME=%h/.config/hermes-ha/local_trees/profiles/telegram2
Environment=HERMES_STRICT_ROUTE_ENABLED=1
Environment=HERMES_STRICT_ROUTE_ACTIVE_PROFILE=telegram2
Environment=HERMES_STRICT_ROUTE_PROFILE=telegram2
Environment=HERMES_STRICT_ROUTE_HOME=%h/.config/hermes-ha/local_trees/profiles/telegram2
Environment=HERMES_STRICT_ROUTE_SERVICE=hermes-gateway-telegram2.service
Environment=HERMES_STRICT_ROUTE_PROVIDER=cursor-gateway
Environment=HERMES_STRICT_ROUTE_MODEL=gpt-5.6-sol
Environment=HERMES_DISABLE_MCP=1
Environment=HERMES_STRICT_ROUTE_BASE_URL=http://127.0.0.1:18080/v1
Environment=HERMES_STRICT_ROUTE_PROBE_TIMEOUT=3
Environment=HERMES_API_TIMEOUT=1860
Environment=HERMES_STREAM_READ_TIMEOUT=1860
ExecStartPre=/usr/bin/env %h/.hermes/hermes-agent/venv/bin/python \
  %h/cursor-gateway/integrations/hermes/cursor-gateway/strict_route_guard.py \
  preflight \
  --profile telegram2 \
  --profiles-file %h/.config/hermes/strict-route-profiles.json \
  --expected-service hermes-gateway-telegram2.service \
  --expected-home %h/.config/hermes-ha/local_trees/profiles/telegram2 \
  --require-protected \
  --check-environment-home
Restart=no
