# infra/hermes-ha/systemd/hermes-ha-gateway-version-sync.service
[Unit]
Description=Mirror verified primary Cursor Gateway revision to vps-band
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=6h
StartLimitBurst=3

[Service]
Type=oneshot
Environment=HERMES_HA_CONFIG=%h/.config/hermes-ha/config.json
Environment=HERMES_HA_NODE_ID=vps-band
ExecStart=/usr/bin/python3 @HERMES_HA_ROOT@/agent/gateway_version_sync.py sync
TimeoutStartSec=2h
Restart=on-failure
RestartSec=30min
UMask=0077
Nice=10
IOSchedulingClass=idle
NoNewPrivileges=true
PrivateTmp=true
