# infra/hermes-ha/systemd/hermes-ha-checkpoint-watchdog.timer
[Unit]
Description=Check Hermes HA gateway checkpoint age every five minutes

[Timer]
OnCalendar=*:3/5
AccuracySec=5s
Persistent=true
Unit=hermes-ha-checkpoint-watchdog.service

[Install]
WantedBy=timers.target
