# infra/hermes-ha/systemd/hermes-ha-evaluate.timer
[Unit]
Description=Timer — Hermes HA evaluate every 2 min (band)

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
Persistent=true
Unit=hermes-ha-evaluate.service

[Install]
WantedBy=timers.target
