# infra/hermes-ha/systemd/hermes-ha-state-checkpoint.timer
[Unit]
Description=Timer — Hermes HA state.db checkpoint every 15 min

[Timer]
OnCalendar=*:1/15
AccuracySec=5s
Persistent=true
Unit=hermes-ha-state-checkpoint.service

[Install]
WantedBy=timers.target
