# infra/hermes-ha/systemd/hermes-ha-gateway-version-sync.timer
[Unit]
Description=Check primary Cursor Gateway revision every 12 hours

[Timer]
OnCalendar=*-*-* 00,12:00:00
Persistent=true
RandomizedDelaySec=45min
FixedRandomDelay=true
AccuracySec=5min
Unit=hermes-ha-gateway-version-sync.service

[Install]
WantedBy=timers.target
